# crontab
* * * * * /home/xc/env/python/logSplit.sh
